# yellowsn0w_0.9.6.deb/System/Library/LaunchDaemons/org.iphone-dev.yellowsn0w.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>org.iphone-dev.yellowsn0w</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/yellowsn0w</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/logs/yellowsn0w.stderr</string>
	<key>StandardOutPath</key>
	<string>/var/logs/yellowsn0w.stdout</string>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
